# Deployments/laife-ai/deploy-bundle/systemd/laife-gpu@.service
[Unit]
Description=Laife GPU service (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Slice=laife-base.slice
User=laife
Group=laife
WorkingDirectory=/opt/laife/apps/laife-ai
EnvironmentFile=/opt/laife/config/services.env
EnvironmentFile=/opt/laife/config/secrets.env
Environment=CUDA_VISIBLE_DEVICES=0
Environment=PYTORCH_CUDA_ALLOC_CONF=max_split_size_mb:512
ExecStart=/opt/laife/bin/launch-%i.sh
Restart=always
RestartSec=10
StandardOutput=append:/opt/laife/logs/laife-ai/%i.log
StandardError=append:/opt/laife/logs/laife-ai/%i.err
LimitNOFILE=65535
TimeoutStartSec=300
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/laife/data /opt/laife/logs /opt/laife/apps /tmp

[Install]
WantedBy=multi-user.target
